# ops/systemd/ip-hot-llm.timer
[Unit]
Description=Drain the IP Hot LLM queue every 3 minutes

[Timer]
OnCalendar=*-*-* *:0/3:00
RandomizedDelaySec=20
Persistent=true
Unit=ip-hot-llm.service

[Install]
WantedBy=timers.target

